# files/etc/profile.d/go-path.sh
export GOPATH=$HOME/go
export GOBIN=$HOME/go/bin
export PATH=$GOPATH/bin:$PATH:/usr/lib/go-1.8/bin
